# source/etc/cron.daily/storecast
#!/bin/bash
storecast -j &>/tmp/storecast.log
