# source/usr/local/emhttp/plugins/dwstorecast/misc/storecast.cron
#!/bin/bash
echo "Generating a new storage forecast..." | logger -t "storecast-cron"
nohup storecast --path /mnt/user --json --normalize &>/tmp/storecast.log &
